# tftp-hpa/rootfs/etc/services.d/tftpd-hba/run
#!/usr/bin/env bashio
networks=$(bashio::config 'networks')
tftproot=$(bashio::config 'tftproot')
pidfile=/tmp/tftp-pidfile.pid

echo "Starting tftpd..."
/usr/sbin/in.tftpd --verbose --foreground --verbose --secure --map-file /etc/tftp.map --address :69 $tftproot}

# /usr/sbin/in.tftpd --verbose --foreground --verbose --secure --pidfile ${pidfile} --map-file /etc/tftp.map --address :69 $tftproot}
# wait "$(cat ${pidfile})"
